# deploy/seromod.service
[Unit]
Description=Seromod Discord Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=seromod
WorkingDirectory=/opt/seromod/seromod
EnvironmentFile=/opt/seromod/seromod/.env
ExecStart=/opt/seromod/seromod/.venv/bin/python bot.py
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ReadWritePaths=/opt/seromod/seromod
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
UMask=027
StandardOutput=journal
StandardError=journal
SyslogIdentifier=seromod

[Install]
WantedBy=multi-user.target
